# offline-crm-v2/deploy/crm-auth.service
[Unit]
Description=PRESSCO21 CRM Auth Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/pressco21/crm-auth/crm_auth_server.py
Restart=always
RestartSec=2
Environment=CRM_AUTH_BIND=127.0.0.1
Environment=CRM_AUTH_PORT=9100
Environment=CRM_AUTH_HTPASSWD_FILE=/etc/nginx/.htpasswd-crm
Environment=CRM_AUTH_SECRET_FILE=/etc/pressco21-crm/auth-secret
Environment=CRM_AUTH_COOKIE_NAME=pressco21_crm_session
Environment=CRM_AUTH_SESSION_TTL=43200
Environment=CRM_AUTH_AUTOMATION_KEY_FILE=/etc/pressco21-crm/automation-key
Environment=CRM_AUTH_AUTOMATION_HEADER_NAME=X-CRM-Automation-Key
Environment=CRM_AUTH_ALLOW_BASIC_AUTH=1

[Install]
WantedBy=multi-user.target
